# scripts/sh.saqoo.canopy.auto-adopt.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Auto-adopt pipeline: checks for Claude Code VSCode extension updates daily at 08:15 local time.

  This is a TEMPLATE with placeholders. Do not install directly.
  Run scripts/install-auto-adopt.sh, which substitutes the placeholders with
  the absolute paths for this checkout and writes the rendered plist to
  ~/Library/LaunchAgents/.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>sh.saqoo.canopy.auto-adopt</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>bash</string>
    <string>__CANOPY_REPO__/scripts/auto-adopt.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/canopy-auto-adopt.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/canopy-auto-adopt.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>
</dict>
</plist>
